# hudson/templates/default/hudson.erb
#!/bin/sh
#
# hudson    Daemon script for Hudson server.
#
# chkconfig: 345 91 10
# description: See https://hudson.dev.java.net/.
# processname: hudson
# pidfile: /var/run/hudson.pid


. /etc/rc.d/init.d/functions
. /etc/sysconfig/network

# Check that networking is up.
[ "${NETWORKING}" = "no" ] && exit 0

export JAVA_HOME=<%= @hudson[:java] %>
export HUDSON_HOME=<%= @hudson[:home] %>
export HUDSON_USER=<%= @hudson[:user] %>

prog="hudson-server"
exec="${HUDSON_HOME}/${prog}"

[ -f /etc/sysconfig/$prog ] && . /etc/sysconfig/$prog

start() {
    [ -x $exec ] || exit 5

    echo -n $"Starting ${prog}: "
    daemon --check ${prog} --user $HUDSON_USER $exec &
    retval=$?
    echo
    return $retval
}

stop() {
    action $"Stopping ${prog}: "
    echo
    killproc ${prog}
    retval=$?
    echo
    return $retval
}


case "$1" in
    start)
        $1
        ;;
    stop)
        $1
        ;;
    status)
        status hudson
        ;;
    restart)
        stop
        start
        ;;
    force-reload)
        stop
        start
        ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart|force-reload}"
        exit 2
        ;;
esac
exit $?
