# hudson/templates/default/hudson.erb
#!/bin/sh
#
# hudson    Daemon script for Hudson server.
#
# chkconfig: 345 91 10
# description: Hudson

### BEGIN INIT INFO
# Provides: hudson
# Required-Start: $local_fs $network
# Required-Stop: $local_fs $network
# Should-Start: $named $time
# Should-Stop: $named $time
# Short-Description: Startup script for hudson 
# Description: See http://hudson-ci.org/
### END INIT INFO

. /etc/rc.d/init.d/functions

JAVA_HOME="<%= @hudson[:java] %>"
HUDSON_USER="<%= @hudson[:user] %>"
HUDSON_HOME="<%= @hudson[:home] %>"

prog="hudson-server"
exec="${HUDSON_HOME}/${prog}"
pidfile="/var/run/hudson.pid"

[ -f /etc/sysconfig/${hudson} ] && . /etc/sysconfig/${hudson}

start() {
    [ -x ${exec} ] || exit 5

    cmd="HUDSON_HOME=${HUDSON_HOME} JAVA_HOME=${JAVA_HOME} ${exec}"
    action $"Starting ${prog}: "
    daemon --pidfile ${pidfile} --user ${HUDSON_USER} "${cmd}"
    `ps -ef | grep hudson.war | grep -v grep | awk '{ print $2 }'` > ${pidfile}
    retval=$?
    echo
    return $retval
}

stop() {
    action $"Stopping ${prog}: "
    killproc -p ${pidfile} ${prog}
    retval=$?
    echo
    return $retval
}


case "$1" in
    start)
        $1
        ;;
    stop)
        $1
        ;;
    status)
        status hudson
        ;;
    restart)
        stop
        start
        ;;
    force-reload)
        stop
        start
        ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart|force-reload}"
        exit 2
        ;;
esac
exit $?
