# hudson/templates/default/hudson.erb
#!/bin/sh
#
# hudson    Daemon script for Hudson server.
#
# chkconfig: 345 91 10
# description: Hudson

### BEGIN INIT INFO
# Provides: hudson
# Required-Start: $local_fs $network
# Required-Stop: $local_fs $network
# Should-Start: $named $time
# Should-Stop: $named $time
# Short-Description: Startup script for hudson 
# Description: See http://hudson-ci.org/
### END INIT INFO

. /etc/rc.d/init.d/functions

JAVA_HOME="<%= @hudson[:java] %>"
JAVA_OPTS="<%= @hudson[:javaopts].join(' ') %>"
HUDSON_USER="<%= @hudson[:user] %>"
HUDSON_HOME="<%= @hudson[:home] %>"
HUDSON_OPTS="<%= @hudson[:opts].join(' ') %>"

prog="hudson"
exec="${HUDSON_HOME}/hudson-server"
pidfile="/var/run/hudson.pid"
cmd="HUDSON_HOME=\"${HUDSON_HOME}\" JAVA_HOME=\"${JAVA_HOME}\" JAVA_OPTS=\"${JAVA_OPTS}\" HUDSON_OPTS=\"${HUDSON_OPTS}\" ${exec}"

[ -f /etc/sysconfig/${prog} ] && . /etc/sysconfig/${prog}

start() {
    [ -x ${exec} ] || exit 5

    action $"Starting ${prog}: "
    daemon --check ${prog} --pidfile ${pidfile} --user ${HUDSON_USER} "${cmd}"
    retval=$?
    if [ $retval -eq 0 ]; then
        pid=`ps -ef | grep hudson.war | grep -v grep | awk '{ print $2 }'`
        echo "${pid}" > ${pidfile}
    fi
    echo
    return $retval
}

stop() {
    action $"Stopping ${prog}: "
    killproc -p ${pidfile} ${prog}
    retval=$?
    echo
    return $retval
}

restart() {
    stop
    start
}

case "$1" in
    start)
        $1
        ;;
    stop)
        $1
        ;;
    status)
        status -p ${pidfile} ${prog}
        ;;
    restart)
        $1
        ;;
    force-reload)
        restart
        ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart|force-reload}"
        exit 2
        ;;
esac
exit $?
