# hudson/templates/default/hudson.erb
#!/bin/sh
#
# hudson    Daemon script for Hudson server.
#
# chkconfig: 345 91 10
# description: Hudson

### BEGIN INIT INFO
# Provides: hudson
# Required-Start: $local_fs $network
# Required-Stop: $local_fs $network
# Should-Start: $named $time
# Should-Stop: $named $time
# Short-Description: Startup script for hudson 
# Description: See http://hudson-ci.org/
### END INIT INFO

. /etc/rc.d/init.d/functions

prog=hudson
pidfile="/var/run/hudson.pid"

JAVA_HOME="<%= @hudson[:java] %>"
JAVA="${JAVA_HOME}/bin/java"
HUDSON_USER="<%= @hudson[:user] %>"
HUDSON_HOME="<%= @hudson[:home] %>"
HUDSON_LOG="${HUDSON_HOME}/hudson.log"
HUDSON="${HUDSON_HOME}/hudson.war"

[ -f /etc/sysconfig/${prog} ] && . /etc/sysconfig/${prog}

start() {
    [ -f ${HUDSON} ] || exit 5

    cmd="HUDSON_HOME=${HUDSON_HOME} JAVA_HOME=${JAVA_HOME} ${JAVA} -jar ${HUDSON} > ${HUDSON_LOG} 2>&1"

    action $"Starting ${prog}: "
    daemon --check ${prog} --pidfile ${pidfile} --user ${HUDSON_USER} "${cmd}"
    retval=$?
    return $retval
}

stop() {
    action $"Stopping ${prog}: "
    killproc -p ${pidfile} ${prog}
    retval=$?
    return $retval
}


case "$1" in
    start)
        $1
        ;;
    stop)
        $1
        ;;
    status)
        status hudson
        ;;
    restart)
        stop
        start
        ;;
    force-reload)
        stop
        start
        ;;
    *)
        echo $"Usage: $0 {start|stop|status|restart|force-reload}"
        exit 2
        ;;
esac
exit $?
